# open_access_scripts/rc.local
# from the man arclight

/home/access/scripts/start_screen_logging.sh